# services/rusty-oxigotchi.service
[Unit]
Description=Rusty Oxigotchi - WiFi capture daemon
After=network.target wlan-keepalive.service
Wants=wlan-keepalive.service

[Service]
Type=simple
ExecStart=/usr/local/bin/rusty-oxigotchi
Environment=RUST_LOG=info
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
